# build/deb/cow/postinst
#!/bin/bash

mountall="/etc/init/mountall.conf"
if test -e $mountall; then
  echo >> $mountall
  echo "pre-start script" >> $mountall
  echo "    exec /etc/init.d/cow.sh start" >> $mountall
  echo "end script" >> $mountall
  name="h"
  if test -e /dev/sda; then
    name="s"
  fi
  tune2fs -c -1 /dev/$name"da1"
  tune2fs -i -1 /dev/$name"da1"
else
  update-rc.d cow.sh start 01 S .
fi

update-rc.d cow_omit_pid defaults
