# build/deb/samba/prerm
#!/bin/bash
service samba stop
sed -i -r 's/SAMBA_ENABLED=true/SAMBA_ENABLED=false/g' /etc/grid_appliance.config
if test -e /etc/samba/smb.conf.orig; then
  mv -f /etc/samba/smb.conf.orig /etc/samba/smb.conf
fi
update-rc.d samba defaults
service samba start
