# build/deb/samba/prerm
#!/bin/bash
smb=samba
if [[ -e /etc/init.d/smbd ]]; then
  smb=smbd
fi

service $smb stop
if test -e /etc/samba/smb.conf.orig; then
  mv -f /etc/samba/smb.conf.orig /etc/samba/smb.conf
fi
update-rc.d $smb defaults
service $smb start
