# build/deb/samba/prerm
#!/bin/bash
service samba stop
rm -f $DIR/etc/samba
if test -e /etc/samba/smb.conf.orig; then
  mv -f /etc/samba/smb.conf.orig /etc/samba/smb.conf
fi
update-rc.d samba defaults
service samba start
