# build/deb/client/prerm
#!/bin/bash
update-rc.d -f startx.sh remove
sed "/eth1/d" -i /etc/network/interfaces
