# build/deb/ssh/postinst
#!/bin/bash
source /etc/grid_appliance.config
update-rc.d -f ssh remove
service ssh stop

if ! test -e /etc/ssh/sshd_config.orig; then
  mv -f /etc/ssh/sshd_config /etc/ssh/sshd_config.orig
fi

sed -i -r 's/SSH_ENABLED=false/SSH_ENABLED=true/g' /etc/grid_appliance.config
service grid_appliance.sh ssh
