# build/deb/ssh/postinst
#!/bin/bash
source /etc/grid_appliance.config
update-rc.d -f ssh remove
service ssh stop

if ! test -e /etc/ssh/sshd_config.orig; then
  mv -f /etc/ssh/sshd_config /etc/ssh/sshd_config.orig
fi

service grid_appliance.sh ssh
