# build/deb/base/postinst
#!/bin/bash
mv /etc/issue /opt/grid_appliance/etc/issue.old

cd /etc/
for i in ipop.vpn.config grid_appliance.config issue; do
  ln -sf ../opt/grid_appliance/etc/$i .
done
cd - &> /dev/null

source /etc/grid_appliance.config
GDIR=$DIR
source /etc/ipop.vpn.config
ln -sf $GDIR/tools/Brunet.Inject.HostActivity.dll $DIR/bin/.

update-rc.d grid_appliance.sh defaults

/etc/init.d/grid_appliance.sh start
