# build/deb/base/postinst
#!/bin/bash
source /etc/ipop.vpn.config
IDIR=$DIR
source /etc/grid_appliance.config

# Store files
mv /etc/issue $DIR/etc/issue.old

# Add files
ln -sf $DIR/etc/issue /etc/.
ln -sf $DIR/etc/apt /etc/apt/apt.conf.d/02ga-upgrades
ln -sf $DIR/scripts/condor_config.sh /etc/condor/.
ln -sf $DIR/tools/Brunet.Inject.HostActivity.dll $IDIR/bin/.

if test -e /etc/init.d/condor; then
  update-rc.d -f condor remove
  /etc/init.d/condor stop
fi

update-rc.d grid_appliance.sh defaults

/etc/init.d/grid_appliance.sh start
