# build/deb/base/prerm
#!/bin/bash
source /etc/ipop.vpn.config
IDIR=$DIR
source /etc/grid_appliance.config

/etc/init.d/grid_appliance.sh stop

# Remove files
rm -f /etc/condor/condor_config.sh \
  $IDIR/bin/Brunet.Inject.HostActivity.dll \
  /etc/apt/apt.conf.d/02ga-upgrades \
  /usr/sbin/condor_kbdd

# Restore files
mv -f $DIR/etc/issue.old /etc/issue

# Update files
sed "/condor_config\.sh/d" -i /etc/condor_config

# Fix rc

if test -e /etc/init.d/condor; then
  update-rc.d condor defaults
  /etc/init.d/condor start
fi

update-rc.d -f grid_appliance.sh remove
