# build/deb/base/prerm
#!/bin/bash
/etc/init.d/grid_appliance.sh stop

mv -f /opt/grid_appliance/etc/issue.old /etc/issue
source etc/ipop.vpn.config
rm -f $DIR/bin/Brunet.Inject.HostActivity.dll

update-rc.d -f grid_appliance.sh remove
